# packages/emacs-echenv/root/_user-home_/profile.d/emacs.sh._template_
## -*- Shell-script -*-
#compiler-settings
cheetahVarStartToken = @
#end compiler-settings

#if @OS_TYPE == 'cygwin' and self.varExists('EMACS_HOME')
    convertenvvar EMACS_HOME 
    pathmunge "${EMACS_HOME}/bin" after
    export EMACS="${EMACS_HOME}/bin/runemacs.exe"
    export EMACS_CLI="${EMACS_HOME}/bin/emacsclient.exe"
#else
    export EMACS_HOME=/usr/bin
    export EMACS="${EMACS_HOME}/emacs"
    export EMACS_CLI="${EMACS_HOME}/emacsclient"
#end if

if [ `expr "$TERM" : "emacs*"` -ne 0 ] || [ `expr "$TERM" : "eterm*"` -ne 0 ] ; then
    export PS1="$ "
fi

function runemacs
{
    "$EMACS" --debug-init
}

export -f runemacs
alias emacs=runemacs

function edit
{
#if @OS_TYPE == 'cygwin' and self.varExists('EMACS_HOME')
    FILE=`cygpath -aw "$*"`
#else
    FILE=$*
#end if
    "$EMACS_CLI" -n "$FILE" > /dev/null 2>&1
    if [ "$?" -ne 0 ] ; then	
	"$EMACS" --debug-init "$FILE" &
    fi
}
export -f edit
export EDITOR="$EMACS_CLI -n"
