# shell/profile.d/python.sh
pathmunge  ${HOME}/.local/bin after
export PYTHONPATH="/home/ech/.local/lib/python2.7/site-packages/:/usr/local/lib/python2.7/dist-packages/:/usr/lib/python2.7:/usr/lib/python2.7/plat-linux2:/usr/lib/python2.7/lib-tk:/usr/lib/python2.7/lib-old:/usr/lib/python2.7/lib-dynload:/usr/local/lib/python2.7/dist-packages"



